# cloudlab_setup/c6220/rapl_service/rapl_log.service
[Unit]
Description=Service to store rapl energy log

[Service]
WorkingDirectory=/users/sanas/experiment-scripts/cloudlab_setup/c6220/rapl-service
ExecStart=/usr/bin/python rapl_log.py

[Install]
WantedBy=multi-user.target